# blfs-servicescripts/services/swat/finish
#!/bin/sh

pidfile=/run/stunnel-swat.pid

pgrep -u root stunnel

rm -rf $pidfile