# trunk/services/dbus/finish
#!/bin/sh
/bin/killall -q /usr/bin/dbus-daemon
rm -rf /run/dbus/system_bus_socket /run/dbus/pid