# trunk/services/kdm/run
#!/bin/sh
sv start /var/service/getty-* || exit 1
exec /usr/bin/kdm -nodaemon