# trunk/rc.local
#!/bin/sh

# This is the rc.local file to load any programs not covered by Runit
# as daemons. For best results, you should use the following sample
# as a method to load a daemon...
#
# Load Timidity++ as Daemon
# if [ -x /usr/bin/timidity ]; then
#     echo "Starting TiMidity++"
#     /usr/bin/timidity -iAD -Os -B2,8
# fi
#
# Please try to avoid scripts from BLFS that utilize init-functions
# from the /lib/lsb directory. File init-functions is not supported by
# Runit-for-LFS.