# trunk/services/avahi-daemon/run
#!/bin/sh
sv check dbus eth0 >/dev/null || exit 1
exec avahi-daemon -s