# trunk/services/wicd/run
#!/bin/sh
sv check dbus eth0 >/dev/null || exit 1
exec /usr/sbin/wicd --no-daemon