# trunk/services/ntpd/finish
#!/bin/sh
kill -9 $(pidof ntpd)
rm /var/run/ntpd.pid
