# trunk/services/ntpd/run
#!/bin/sh
exec /usr/sbin/ntpd  --nofork -g -u ntp:ntp