# trunk/services/ntpd/run
#!/bin/sh
exec ntpd --nofork -g -u ntp:ntp