# trunk/services/sshd/run
#!/bin/sh
sv check eth0 >/dev/null || exit 1
/usr/bin/ssh-keygen -A # Will generate host keys if they don't already exist
exec /usr/sbin/sshd -D