# trunk/services/sshd/run
#!/bin/sh

# sshd will crash if you use the BLFS example to create
# the keys. Let the Runit run file create them for you.

/usr/bin/ssh-keygen -A

exec /usr/sbin/sshd -D 2>&1
