# trunk/blfs-servicescripts/services/winbindd/finish
#!/bin/sh

PIDFILE=/var/run/winbindd.pid

pgrep -u root winbindd

rm -rf ${PIDFILE}