# trunk/services/fcron/run
#!/bin/sh
# Begin fcron service run file
exec fcron -f
# End of fcron service run file
