# trunk/services/gpm/finish
#!/bin/sh
kill -9 $(pidof gpm) &>/dev/null
rm -f /run/gpm.pid > &>/dev/null