# trunk/services/gpm/finish
#!/bin/sh
rm -f /run/gpm.pid &>/dev/null