# trunk/services/gpm/finish
#!/bin/sh
# Begin gpm service finish file
#Stop the mouse daemon.

kill -9 $(pidof gpm) &>/dev/null
rm -f /run/gpm.pid > &>/dev/null

#End /etc/sv/gpm/run