# blfs-servicescripts/services/avahi-daemon/run
#!/bin/sh
sv check /etc/runit/runsvdir/multi/dbus >/dev/null || exit 1
exec avahi-daemon -s