# trunk/services/smartd/finish
#!/bin/sh
kill -9 $(pidof /usr/sbin/smartd)
rm -f /run/smartd.pid
