# trunk/services/smartd/run
#!/bin/sh
exec /usr/sbin/smartd -n  -p /run/smartd.pid  -q never 2>&1
