# blfs-servicescripts/services/dhcpd/run
#!/bin/sh
sv check eth0 >/dev/null || exit 1
exec 2>&1
exec dhcpd -f -d -cf ./config eth0