# trunk/services/agetty-ttyS0/run
#!/bin/sh
exec chpst -P /sbin/agetty -8 -s 38400 --noclear tty1 linux