# trunk/services/cupsd/finish
#!/bin/sh
# Stop the printer daemon.
/bin/killall -q /usr/sbin/cupsd
# End /etc/sv/cups/finish