# trunk/services/cupsd/finish
#!/bin/sh
/bin/killall -q /usr/sbin/cupsd